# scripts/systemd/leviculum-ci-nightly.timer
[Unit]
Description=Leviculum CI nightly schedule

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
