# scripts/systemd/leviculum-ci-nightly.service
[Unit]
Description=Leviculum CI Tier 3 (nightly)

[Service]
Type=oneshot
WorkingDirectory=%h/coding/libreticulum
ExecStart=/bin/bash %h/coding/libreticulum/scripts/run-tier3.sh
