# packaging/debian/lnsd.service
[Unit]
Description=Leviculum / Reticulum Network Stack Daemon
Documentation=https://codeberg.org/Lew_Palm/leviculum
After=network.target
Wants=network.target

[Service]
Type=simple
User=leviculum
Group=leviculum
ExecStart=/usr/bin/lnsd -c /etc/reticulum
Restart=on-failure
RestartSec=5

# Socket + storage files created by the daemon are owned leviculum:leviculum
# with mode 0660 / dir 0770 — only members of the `leviculum` group can
# reach the shared-instance socket.
UMask=0007

# Hardening: read-only filesystem except the config/state directory.
ProtectSystem=strict
ReadWritePaths=/etc/reticulum
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
